# ai/.config/systemd/user/opencode-web.service
[Unit]
Description=Opencode Web Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/bin/opencode-web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
